# base-files/eatsa-user/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Manually created to let the eatsa user automatically login.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin eatsa --noclear %I $TERM
