# base-files/supervisor-scripts/etc/supervisor/conf.d/sup-smartshelf.conf
# /etc/supervisor/conf.d/
[program:smartshelf-server]
command=nodejs /home/eatsa/Documents/Smartshelf/master.js

autorestart=true

stdout_logfile=AUTO
redirect_stderr=true
